# workspace/crontab/laradock
* * * * * laradock php /var/www/artisan schedule:run >> /dev/null 2>&1
0 */1 * * * laradock python3 /var/www/EasemobSDK/start.py >> /var/www/EasemobSDK/log.log 2>&1
